# openrc/servicio-tvheadend
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
    need net
    need localmount
    need nfsmount
    after sshd
}

start() {
	ebegin "Arrancando servicio Tvheadend"

        su luis -c "cd /Apps/docker/servicio-tvheadend && /usr/local/bin/fig up -d"

        eend $?
}

stop() {
	ebegin "Parando servicio Tvheadend"

        su luis -c "cd /Apps/docker/servicio-tvheadend &&  /usr/local/bin/fig stop"

	eend $?
}

status() {

        su luis -c "cd /Apps/docker/servicio-tvheadend &&  /usr/local/bin/fig ps"

	eend $?
}
